# ops/systemd/adaad-forensic-retention.timer
[Unit]
Description=Run ADAAD forensic retention enforcement daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=adaad-forensic-retention.service

[Install]
WantedBy=timers.target
